# modules/irc_slack/templates/irc-slack.service.erb
[Unit]
Description=IRC-to-Slack gateway
After=network.target
Documentation=<%= @soucre %>

[Service]
ExecStart=/usr/local/bin/irc-slack <% if @debug -%>-D <% end -%>-H <%= @listen %> -p <%= @port %> -s <%= @server_name %>
Restart=always

[Install]
WantedBy=multi-user.target
